# addons_local/xcomfort_clean/xcomfort_go2/rootfs/etc/service/xcomfortd/run
#!/bin/sh
exec /bin/busybox sh -c 'echo "[xcomfort] minimal legacy run"; sleep 3600'
