# addons_local/xcomfort_clean/rootfs/etc/services.d/xcomfortd/run
#!/bin/sh
exec /bin/busybox sh -c 'echo "[xcomfort] CLEAN minimal run"; sleep 3600'
