# addons_local/xcomfort_go/rootfs/etc/services.d/xcomfortd/run
#!/bin/sh
exec /bin/busybox sh -c 'echo "[xcomfort] GO minimal run"; sleep 3600'
